# kernel-work/morse-kernel-staging-cm4/etc/udev/rules.d/80-ethernet-autodetect.rules
# Trigger on carrier changes for ethernet interface
SUBSYSTEM=="net", KERNEL=="end0", ACTION=="change", \
  RUN+="/bin/systemctl start ethernet-autodetect.service"
